# homeserver/launchd/com.cubinghub-backup.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cubinghub.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd "__REPO_DIR__" &amp;&amp; ./homeserver/scripts/backup-home-server.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME_DIR__/Library/Logs/cubing-hub-backup.out.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME_DIR__/Library/Logs/cubing-hub-backup.err.log</string>

  <key>WorkingDirectory</key>
  <string>__REPO_DIR__</string>
</dict>
</plist>
